# services/spot-mcp-server/memory-janitor.service
[Unit]
Description=Spot Memory Janitor - Autonomous memory maintenance
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/docker run --rm --network="host" -e QDRANT_URL=http://localhost:6333 -e COLLECTION_NAME=default-collection spot-mcp-server:arm64 python -m memory_janitor
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
